# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; ── php-fpm ──────────────────────────────────────────────────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Laravel scheduler (runs every minute via sleep loop) ─────────────────────
; The scheduler command is lightweight — it wakes up, checks for due tasks,
; and exits. Supervisor restarts it immediately, giving us a ~1-minute cadence.
[program:scheduler]
command=/bin/bash -c "while true; do php /var/www/artisan schedule:run --no-interaction; sleep 60; done"
autostart=true
autorestart=true
priority=30
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
